# .config/systemd/user/batsignal.service
[Unit]
Description=batsignal

[Service]
Type=simple
ExecStart=/usr/bin/batsignal -w 20 -c 10 -d 5

[Install]
WantedBy=default.target
