# .config/systemd/user/lowbatteryalert.timer
[Unit]
Description=Notify when battery is below 20%

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
