# .config/systemd/user/onedrive.service
[Unit]
Description=Sync OneDrive to local directory using rclone
After=network-online.service
Wants=network-online.service

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync --check-access -v onedrive: %h/onedrive
ExecStopPost=/usr/bin/sh -c 'if [ "$$SERVICE_RESULT" != "success" ]; then notify-send "OneDrive sync failed"; fi'

[Install]
WantedBy=default.target
