# .config/systemd/user/onedrive.path
[Path]
PathModified=%h/onedrive

[Install]
WantedBy=default.target
