# .config/systemd/user/restic.service
[Unit]
Description=Restic backup to Backblaze B2
After=network-online.service
Wants=network-online.service

[Service]
Type=oneshot
EnvironmentFile=%h/.restic/env
ExecStart=/usr/bin/restic backup --verbose --files-from %h/.restic/files
ExecStartPost=/usr/bin/restic forget --keep-last 3 --prune

[Install]
WantedBy=default.target
