# .config/systemd/user/restic.timer
[Unit]
Description=Run restic backup weekly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
